# scripts/macmini-translate-worker/com.goldenfocus.dalat-translate-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.goldenfocus.dalat-translate-worker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/theoutsider/dalat-app/scripts/macmini-translate-worker/run.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- Pinned at deploy time to the moment the last qwen3 writer was
         killed. REDO_BEFORE-REPLACED-AT-DEPLOY -->
    <key>REDO_BEFORE</key>
    <string>REPLACE_AT_DEPLOY</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/Users/theoutsider/Library/Logs/dalat-translate-worker.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/theoutsider/Library/Logs/dalat-translate-worker.log</string>
</dict>
</plist>
